# scripts/repository/inst-min-con-xtra-macchanger/macchanger-udev.rules
SUBSYSTEM=="net", ENV{DEVTYPE}=="wlan", RUN+="/usr/bin/macchanger -e $env{INTERFACE}"

